# platform/initrd/common/etc/inittab
# /etc/inittab

::sysinit:/etc/init.d/rc.startup

tty1::respawn:/sbin/getty -nl /boot/system-status 38400 tty1
tty6::askfirst:/sbin/getty -nl /boot/recovery-console 38400 tty6

::restart:/etc/init.d/rc.shutdown
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::shutdown:/etc/init.d/rc.shutdown

